# bardinux-unattended-upgrades/gcs/newfiles_skel/etc/init/unattended-upgrades.conf
description		"Run unattended upgrades on network up"
author			"Enrique Hernández Bello <quique@osl.ull.es>"

start on net-device-up
exec /etc/cron.daily/apt
