# bardinux-default-settings/gcs/newfiles_skel/usr/share/autostart/create-contact-dir.desktop
[Desktop Entry]
Exec=mkdir -p $HOME/.local/share/contacts
