# bardinux-sigaull-default-settings/gcs/newfiles_skel/etc/cron.d/shutdown
# m h  dom mon dow   command
5 22 * * * 	root	/sbin/init 0
