# bardinux-sigaull-default-settings/gcs/newfiles_skel/etc/cron.d/shutdown
# m h  dom mon dow   command
0 */1 * * * 	root	test -z "$(who)" && test -z "$(rawtmp -b 0:15 | tail -n +2)" && /sbin/shutdown -h now
