# bardinux-sigaull-default-settings/gcs/newfiles_skel/etc/cron.d/shutdown
# m h  dom mon dow   command
0 */2 * * * 	root	test $(who | wc -l) -eq 0 && /sbin/init 0
