# bardinux-sigaull-default-settings/gcs/newfiles_skel/etc/cron.d/shutdown
# m h  dom mon dow   command
0 */1 * * * 	root	test -n $(who) || /sbin/shutdown -h now
