# bardinux-sigaull-default-settings/gcs/newfiles_skel/etc/cron.d/silence
# m h  dom mon dow   command
* * * * *	root	/usr/bin/amixer sset "Master" 0
* * * * *	root	/usr/bin/amixer sset "Master Mono" 0
